# deploy/vk-ads-sync.timer
[Unit]
Description=Run VK Ads synchronization every day at 01:00 Moscow time

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true
Unit=vk-ads-sync.service

[Install]
WantedBy=timers.target
